# gotify/rootfs/etc/s6-overlay/s6-rc.d/gotify/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Gotify
# Runs Gotify
# ==============================================================================
bashio::log.info 'Starting Gotify...'

# Run Gotify
exec /opt/gotify/gotify
